# conf/supervisor.conf
# -*- conf -*-

[program:mapi]
command=/srv/sites/mapi/conf/startup
user=ebower
group=web
autostart=true
autorestart=true
stderr_logfile=/srv/log/mapi/supervisor_err.log
